# package/carbonio-docs-connector.service
# SPDX-FileCopyrightText: 2022 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

[Unit]
Description=Carbonio docs connector daemon
Wants=network.target
Requires=network-online.target
After=network-online.target
PartOf=carbonio.target

[Service]
Type=simple
ExecStart=/usr/bin/carbonio-docs-connector
User=carbonio-docs-connector
Group=carbonio-docs-connector
Restart=on-failure
RestartSec=15
TimeoutSec=60
TimeoutStopSec=120


[Install]
WantedBy=multi-user.target
