# package/carbonio-docs-connector-sidecar.service
# SPDX-FileCopyrightText: 2022 Zextras <https://www.zextras.com>
#
# SPDX-License-Identifier: AGPL-3.0-only

[Unit]
Description=Carbonio docs connector sidecar proxy
Documentation=https://docs.zextras.com/zextras-suite-documentation/latest/
Requires=network-online.target
After=network-online.target
PartOf=carbonio.target

[Service]
Type=simple
ExecStart=/usr/bin/consul connect envoy \
    -token-file /etc/carbonio/docs-connector/service-discover/token \
    -admin-bind localhost:0 \
    -sidecar-for carbonio-docs-connector
Restart=on-failure
RestartSec=15
User=carbonio-docs-connector
KillMode=process
KillSignal=SIGKILL
LimitNOFILE=65536
TimeoutSec=120
TimeoutStopSec=120

[Install]
WantedBy=carbonio.target
